# roles/podman/templates/podman-auto-update.timer
# {{ ansible_managed }}
[Timer]
OnCalendar=
OnCalendar={{ podman__auto_update_service_schedule }}
