# roles/rpm-ostree/templates/rpm-ostreed-automatic.timer
# {{ ansible_managed }}
[Timer]
OnBootSec=
OnUnitInactiveSec=
OnCalendar=
OnCalendar={{ rpm_ostree__upgrade_schedule }}
